# systemd/coding-agent-poll@.timer
[Unit]
Description=Trigger coding-agent-poll periodically (instance: %i)

[Timer]
Unit=coding-agent-poll@%i.service
OnBootSec=30s
# 注：这里不能用变量，systemd timer 不展开 EnvironmentFile。如果要改 interval，
# 编辑此模板或用 systemctl --user edit coding-agent-poll@%i.timer 做 drop-in。
OnUnitActiveSec=60s
Persistent=true

[Install]
WantedBy=timers.target
